# nodeapp/supervisord.conf
[supervisord]
nodaemon=true

[program:blobapi]
command=/usr/bin/npm start
user=root